# systemd/brele-daily-operations.timer
[Unit]
Description=Run Brele Operations daily jobs at 23:30 server time

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=true
Unit=brele-daily-operations.service

[Install]
WantedBy=timers.target
